# scripts/ternakos-scraper.service
[Unit]
Description=TernakOS Harga Scraper
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/ternakos
EnvironmentFile=/home/ubuntu/ternakos/.env
ExecStart=/home/ubuntu/ternakos/.venv/bin/python scripts/ternakos_harga_scraper.py --daemon
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/ubuntu/ternakos/scripts/ternakos_scraper.log
StandardError=append:/home/ubuntu/ternakos/scripts/ternakos_scraper.log

[Install]
WantedBy=multi-user.target
